# gram/jobmanager/auditing/source/globus-gram-audit.cron
#! /bin/sh

prefix="${GLOBUS_LOCATION-/usr/local/globus-6}"
exec_prefix="${prefix}"

if [ -x ${exec_prefix}/sbin/globus-gram-audit ] ; then
${exec_prefix}/sbin/globus-gram-audit --quiet
fi

exit 0
